# system/systemd/iccp-ros-monitor.service.d/10-robot-namespace.conf
[Service]
Environment=ICCP_ROBOT_NAMESPACE=tracer3
